# config/tdl/crontab
# Download videos from Saved Messages daily at 3am (dedup makes reruns safe)
0 3 * * * /bin/sh /etc/tdl/download-saved.sh >> /downloads/tdl-download.log 2>&1
