# systemd/potatostack-swap.service
[Unit]
Description=PotatoStack Swap File Manager
Documentation=https://github.com/fishingpvalues/potatostack
Before=docker.service
After=mnt-seconddrive.mount
Requires=mnt-seconddrive.mount
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/ensure-potatostack-swap.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
