# systemd/potatostack.service
[Unit]
Description=PotatoStack Docker Compose Stack
Documentation=https://github.com/fishingpvalues/potatostack
Requires=docker.service potatostack-swap.service
After=docker.service potatostack-swap.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/USER/potatostack
User=USER
Group=USER

# Pull latest images quietly before starting (optional)
ExecStartPre=-/usr/bin/docker-compose pull --quiet --ignore-pull-failures

# Start the stack
ExecStart=/usr/bin/docker-compose up -d --remove-orphans

# Stop the stack
ExecStop=/usr/bin/docker-compose down

# Reload configuration
ExecReload=/usr/bin/docker-compose pull --quiet
ExecReload=/usr/bin/docker-compose up -d --remove-orphans

# Timeout settings
TimeoutStartSec=600
TimeoutStopSec=60

# Restart policy
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
